# todo/configuration/conf/crontabs/root
# Detect if there is a minieap process every minute
* * * * * /usr/bin/detect_minieap.sh

# Reboot device on 6:30am
30 6 * * * reboot

# Update googlehosts on 12:00am
0 12 * * * /usr/bin/update_host.sh

# Send traffic counter
30 23 * * * /root/traffic.sh
